# .claude/skills/run-model-eval/com.pinchy.eval-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pinchy.eval-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/.pinchy-eval-watchdog/watchdog.sh</string>
    </array>

    <!-- Run at load and every 15 minutes thereafter. -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>900</integer>

    <!-- CRITICAL: do not kill the watchdog's descendant process group when the
         watchdog script exits — the sweep it launches (nohup + disowned) must
         keep running after each 15-minute watchdog invocation returns. -->
    <key>AbandonProcessGroup</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.pinchy-eval-watchdog/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.pinchy-eval-watchdog/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
